# etc/cron.d/dork-detector
# Detector for Linux/Cdorked.A

0 * * * * root /usr/bin/dork-detector

